# setup/template.service
[Unit]
Description=

[Service]
Type=simple
ExecStart=
User=

[Install]
WantedBy=multi-user.target
